# lego_walk.service
[Unit]
Description=Lego Walk Tracker
After=network.target

[Service]
Type=simple
User=christiangoff
WorkingDirectory=/home/christiangoff/lego_walk
Environment="REBRICKABLE_API_KEY=a623b0e5ccd1c8dbae7ddf13cc1f8fe8"
Environment="SMTP_HOST=smtp.gmail.com"
Environment="SMTP_USER=christian.goff@gmail.com"
Environment="SMTP_PASSWORD=your-app-password-here"
ExecStartPre=/home/christiangoff/lego_walk/venv/bin/python3 /home/christiangoff/lego_walk/migrate.py
ExecStart=/home/christiangoff/lego_walk/venv/bin/python3 /home/christiangoff/lego_walk/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
